# packaging/ubuntu/postinst
echo "Starting agent"
if [ $(cat /proc/1/comm) = init ]
then
    start clunk80-ssm-agent || true
elif [ $(cat /proc/1/comm) = systemd ]
then
    systemctl enable clunk80-ssm-agent
    systemctl start clunk80-ssm-agent
    systemctl daemon-reload
fi
