# packaging/ubuntu/postinst
echo "Starting agent"
if [ $(cat /proc/1/comm) = init ]
then
    start bzero-ssm-agent || true
elif [ $(cat /proc/1/comm) = systemd ]
then
    systemctl enable bzero-ssm-agent
    systemctl start bzero-ssm-agent
    systemctl daemon-reload
fi
