# packaging/ubuntu/prerm
echo "Stopping agent"
if [ $(cat /proc/1/comm) = init ]
then
    stop bzero-ssm-agent || true
elif [ $(cat /proc/1/comm) = systemd ]
then
    systemctl stop bzero-ssm-agent
    systemctl disable bzero-ssm-agent
    systemctl daemon-reload
fi
